# home/.local/share/services/wob/run
#!/bin/bash

set -euo pipefail

SOCKET="${XDG_RUNTIME_DIR}/wob.socket"
rm -f "$SOCKET"

mkfifo "$SOCKET"
exec tail -f "$SOCKET" | /usr/bin/wob
