# bootstrap/rc.local
#!/bin/bash

mkdir -p /run/user/cartolari
chown cartolari:cartolari /run/user/cartolari

exit 0
