# bootstrap/rc.local
#!/bin/bash

mkdir -p /run/user/cartolari
chown cartolari:cartolari /run/user/cartolari
mkdir -p /home/cartolari/code
chown cartolari:cartolari /home/cartolari/code
cryptsetup luksOpen --key-file /home/cartolari/.code-keyfile /var/lib/crouton-shared/docker/code.img code
mount /dev/mapper/code /home/cartolari/code

su -c '/usr/bin/ssh-agent -D -a /run/user/cartolari/ssh-agent.socket' cartolari &

exit 0
