# home/.local/share/services/avizo/run
#!/bin/bash

exec /usr/local/bin/avizo-service
