# bootstrap/nix-kvm/upstart-services/dev-vm.conf
description "Run a development virtual machine using QEMU with KVM"
author      "bruno.cartolari@gmail.com"

start on stopping mount-nix
stop on stopping ui
respawn limit 10 60
post-stop exec sleep 5

env LD_LIBRARY_PATH=

exec /home/chronos/user/.nix-profile/bin/qemu-system-x86_64 \
  -chardev socket,path=/tmp/qga.sock,server,nowait,id=qga0 \
  -cpu host \
  -device virtio-balloon \
  -device virtio-net-pci,netdev=net0,mac=DE:AD:BE:EF:EB:4C \
  -device virtio-serial \
  -device virtserialport,chardev=qga0,name=org.qemu.guest_agent.0 \
  -display vnc=:5 \
  -drive file=/home/chronos/user/vms/dev-vm.qcow2,if=virtio \
  -enable-kvm \
  -m 8192m \
  -netdev tap,id=net0,script=/etc/qemu-ifup,downscript=no \
  -qmp unix:/home/chronos/user/vms/qmp.sock,server,nowait \
  -smp sockets=1,cpus=4,cores=2 \
  -vga virtio
