# bootstrap/nix-kvm/upstart-services/dev-vm-forward-port-443.conf
description "Forward port 443 to the DEV VM"
author      "bruno.cartolari@gmail.com"

start on starting network-services
respawn
respawn limit 10 5

env LD_LIBRARY_PATH=

exec /home/chronos/user/.nix-profile/bin/socat TCP-LISTEN:443,fork TCP:10.250.250.15:443
