# bootstrap/nix-kvm/upstart-services/dev-vm-forward-ports.conf
description "Forward ports to the DEV VM"
author      "bruno.cartolari@gmail.com"

start on starting network-services
task

env LD_LIBRARY_PATH=

script
  /home/chronos/user/.nix-profile/bin/socat TCP-LISTEN:80,fork TCP:10.0.2.15:80 &
  /home/chronos/user/.nix-profile/bin/socat TCP-LISTEN:443,fork TCP:10.0.2.15:443 &
end script
