# bootstrap/nix-kvm/upstart-services/mount-nix.conf
description "Mount /nix"
author      "bruno.cartolari@gmail.com"

start on start-user-session
task

script
  mount -o remount,rw /

  if ! mountpoint /nix &> /dev/null; then
    mount --bind /home/chronos/user/nix /nix
  fi
  mount -o remount,exec /nix
end script
