# bootstrap/nix-kvm/upstart-services/unmount-nix.conf
description "Unmount /nix"
author      "bruno.cartolari@gmail.com"

start on stopping ui
task

script
  if mountpoint /nix &> /dev/null; then
    umount /nix
  fi
end script
