# upstart-services/enable-hdmi-audio.conf
description "Start the crouton clipboard service"
author      "bruno.cartolari@gmail.com"

start on start-user-session
task

exec amixer -c 0 sset IEC958,0 unmute
