# home/.local/share/services/inactive-windows-transparency/run
/usr/share/sway-contrib/inactive-windows-transparency.py