# home/.local/share/services/wob-volume-subscriber/run
#!/bin/sh

exec ~/.local/bin/wob-subscribe-to-volume-changes
