# home/.config/systemd/user/kanshi.service
[Unit]
Description=Sway layout manager

[Service]
ExecStart=/usr/bin/kanshi
Restart=always
Type=simple

[Install]
WantedBy=default.target
