# home/.config/systemd/user/wob.socket
[Unit]
Description=WOB (Wayland Overlay Bar) Socket

[Socket]
ListenFIFO=%t/wob.socket
Service=wob.service
