# home/.config/systemd/user/wob.service
[Unit]
Description=WOB (Wayland Overlay Bar)

[Service]
ExecStart=/usr/bin/wob --verbose
Sockets=wob.socket
StandardInput=socket
Type=simple
